# roles/teiserver/files/acme.service
[Unit]
Description=Update TLS certificates

[Service]
ExecStart=acmetool --xlog.journal --xlog.severity=DEBUG --batch reconcile
