# roles/teiserver/files/geoipupdate.service
[Unit]
Description=Update GeoIP database

[Service]
ExecStart=/usr/local/bin/geoipupdate
