# roles/teiserver/files/tei-backup.service
[Unit]
Description=Take backup

[Service]
Nice=10
CPUSchedulingPolicy=batch
IOSchedulingPriority=6
ExecStart=/usr/local/bin/tei-backup-take
