# roles/teiserver/files/update-kicker-keys.service
[Unit]
Description=Update SSH authorized keys for kicker user

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update-kicker-keys.sh
